# etc/systemd/system/navcoin-repair.service
[Unit]
Description=navcoin-repair
After=network.target
After=network-online.target
After=time-sync.target

[Service]
User=odroid
Type=notify
WorkingDirectory=/home/odroid/
ExecStart=/usr/local/bin/navcoind -zapwallettxes=1
ExecStop=/bin/kill $MAINPID
ExecReload=/bin/kill -HUP $MAINPID
TimeoutStopSec=300
StandardOutput=append:/var/log/navcoin-core.log
StandardError=append:/var/log/navcoin-core.log

[Install]
WantedBy=multi-user.target
